# crontab
# runs a scripts to update received messages
* * * * * root /root/api/update.sh